# apps/test-brain/deploy/personal-brain.service
[Unit]
Description=Personal Brain MCP Server
Documentation=https://github.com/yeehaa123/brains
After=network.target

[Service]
Type=simple
ExecStart=/opt/personal-brain/brain
Restart=always
RestartSec=10

# User and permissions
User=personal-brain
Group=personal-brain

# Environment
Environment="NODE_ENV=production"
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/root/.bun/bin"
EnvironmentFile=/opt/personal-brain/.env

# Working directory
WorkingDirectory=/opt/personal-brain

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/personal-brain/data /opt/personal-brain/brain-repo /opt/personal-brain/website /opt/personal-brain/.matrix-storage /opt/personal-brain/node_modules /opt/personal-brain/.bun /opt/personal-brain/.npm
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

# Resource limits
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=personal-brain

[Install]
WantedBy=multi-user.target